# contrail-controller/files/cron.d/nagios-contrail-controller-db-consistency-check-dump
# DB dump performed inline because of bug with db_manage.py: https://jira.tungsten.io/browse/TFB-1620
53 * * * * root cd /var/lib/nagios && docker exec -t configapi_api_1 python2 -c 'from vnc_cfg_api_server.db_manage import _parse_args, db_check; args, api_args = _parse_args("check");  verb = args.operation; db_check(args, api_args)' > contrail-controller-db-consistency-check.dump

