# contrail-controller/files/cron.d/nagios-contrail-controller-rabbitmq-stats
*/5 * * * * root /usr/local/lib/nagios/plugins/collect_rabbitmq_stats.sh
